# backendapptutor.service
[Unit]
Description=Backend app tutor project

[Service]
User=root
Restart=on-failure
ExecStart=/var/lib/jenkins/workspace/Backend-nodejs-pipeline/systemd_script.sh

[Install]
WantedBy=multi-user.target
